# units/http-daemon.service
[Unit]
Description=HTTP Daemon for Trove
After=network.target

[Service]
ExecStart=/usr/sbin/httpd -u git:git -h /var/www/htdocs -f
Restart=always
