# units/morph-cache-server.service
[Unit]
Description=Morph cache server
After=network.target

[Service]
ExecStart=/usr/bin/morph-cache-server
